# home-panel/rootfs/etc/services.d/nginx/run
#!/command/with-contenv bashio
# ==============================================================================
# Home Assistant Community Add-on: Home Panel
# Runs the NGINX daemon
# ==============================================================================
# Wait for Home Panel to become available
bashio::net.wait_for 8234

bashio::log.info "Starting NGINX..."
exec nginx
