# home-panel/rootfs/etc/services.d/nginx/run
#!/usr/bin/with-contenv bash
# ==============================================================================
# Community Hass.io Add-ons: Home Panel
# Starts server to serve the webapp
# ==============================================================================
# shellcheck disable=SC1091
source /usr/lib/hassio-addons/base.sh

hass.log.info "Run nginx server.."

nginx -g "daemon off;"
