# home-panel/rootfs/etc/services.d/nginx/run
#!/usr/bin/with-contenv bash
# ==============================================================================
# Community Hass.io Add-ons: Home Panel
# Runs the NGINX daemon
# ==============================================================================
# shellcheck disable=SC1091
source /usr/lib/hassio-addons/base.sh

hass.log.info "Starting NGINX for Riot web..."
if hass.config.true 'ssl'; then
    options='/etc/nginx/nginx-ssl.conf'
else
    options='/etc/nginx/nginx.conf'
fi

exec nginx -c "$options"
exec sleep 10
