# home-panel/rootfs/etc/services.d/nginx/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Community Hass.io Add-ons: Home Panel
# Runs the NGINX daemon
# ==============================================================================
# Wait for the API to start
bashio::net.wait_for 3234

# Disable HA Authentication if front door is open
if bashio::config.true 'leave_front_door_open'; then
    export DISABLE_HA_AUTHENTICATION=true
fi

bashio::log.info "Starting NGINX..."
exec nginx
