# home-panel/rootfs/etc/services.d/nginx/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Community Hass.io Add-ons: Home Panel
# Runs the NGINX daemon
# ==============================================================================
bashio::log.info "Starting NGINX..."

if bashio::config.true 'ssl'; then
    options='/etc/nginx/nginx-ssl.conf'
else
    options='/etc/nginx/nginx.conf'
fi

exec nginx -c "$options"
