# home-panel/rootfs/etc/services.d/home-panel/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Community Hass.io Add-ons: Home Panel
# Starts server to serve the app
# ==============================================================================
bashio::log.info "Run Home Panel.."

cd /opt/panel \
    || bashio::exit.nok "Could not change into app directory"

SUPPRESS_URL="true"

export SUPPRESS_URL

exec node lib/
