# home-panel/rootfs/etc/services.d/api/run
#!/usr/bin/with-contenv bash
# ==============================================================================
# Community Hass.io Add-ons: Home Panel
# Starts server to serve the api
# ==============================================================================
# shellcheck disable=SC1091
source /usr/lib/hassio-addons/base.sh

hass.log.info "Run API.."

# shellcheck disable=SC2164
cd /usr/src/api

node index.js
