# home-panel/rootfs/etc/services.d/api/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Community Hass.io Add-ons: Home Panel
# Starts server to serve the API
# ==============================================================================
bashio::log.info "Run API.."

cd /opt/api \
    || bashio::exit.nok "Could not change into API directory"

LOG_LEVEL=$(bashio::config 'log_level')

if bashio::config.true 'ssl'; then
    SSL_PATH_CERT="/ssl/$(bashio::config 'certfile')"
    SSL_PATH_KEY="/ssl/$(bashio::config 'keyfile')"
fi

export LOG_LEVEL
export SSL_PATH_CERT
export SSL_PATH_KEY

exec node src/
